# 56-systemd/journald_logs_kernel_log/dummy-journald.service
[Unit]
Description=Dummy Journald Service for Kernel Logs
After=network.target

[Service]
ExecStart=/usr/local/bin/dummy-journald
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
