# 56-systemd/journald_add_new_service/my_service.service
[Unit]
Description=My Example Service
After=network.target

[Service]
ExecStart=/usr/local/bin/my_service
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
